# automagik/cli/templates/automagik.service
[Unit]
Description=AutoMagik Service
After=network.target postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=%USER%
Group=%USER%
WorkingDirectory=%WORKDIR%
Environment=PYTHONUNBUFFERED=1
Environment=AUTOMAGIK_DEBUG=1
Environment=AUTOMAGIK_LOG_LEVEL=DEBUG
EnvironmentFile=%WORKDIR%/.env

# Start both the scheduler and Celery worker
ExecStart=/bin/bash -c '\
    %VENV_PATH%/bin/celery -A core.celery worker --loglevel=info & \
    %VENV_PATH%/bin/automagik run start --daemon'

# Ensure both processes are stopped
ExecStop=/bin/bash -c 'pkill -f "celery worker" || true'

StandardOutput=journal
StandardError=journal
SyslogIdentifier=automagik
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
