# automagik/api/automagik-api.service
[Unit]
Description=AutoMagik API Service
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/automagik
Environment=PYTHONPATH=/root/automagik
Environment=LANGFLOW_API_URL=http://localhost:7860
Environment=TIMEZONE=UTC
EnvironmentFile=/root/automagik/.env
ExecStart=/root/automagik/.venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
